# diaspora.service
[Unit]
Descript=Diaspora Server

[Service]
Type=simple
ExecStart=/opt/diaspora
WorkingDirectory="/home/adm001di"
User="adm001di"


