# diaspora.service
[Unit]
Description=Diaspora socail network (application server)
Wants=redis-server.service
Wants=mysql.service

[Service]
User=adm001di
WorkingDirectory=/home/adm001di/diaspora
Environment=RAILS_ENV=production DB=mysql
ExecStart=/bin/bash -lc "true; /opt/./diaspora"
CPUAccounting=true
MemoryAccounting=true
BlockIOAccounting=true

[Install]
WantedBy=multi-user.target
