# diaspora.service
[Unit]
Descript=Diaspora Server

[Service]
Type=simple
ExecStart=/opt/diaspora
WorkingDirectory="~"
User="adm001di"


